# radar_detector.service
[Unit]
Description=Automatic Radar Detection System
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/radar-detector
ExecStart=/usr/bin/python3 /home/pi/radar-detector/radar_detector_auto.py /home/pi/radar-detector/radars.csv
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target