# radar_detector.service
[Unit]
Description=Automatic Radar Detection System
After=network.target
Wants=network.target

[Service]
Type=simple
User=noagiannone
Group=noagiannone
ExecStart=/usr/bin/python3 /home/noagiannone/raspi-radar-detector/radar_detector_auto.py /home/noagiannone/raspi-radar-detector/radars.csv
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target